# configs/tailscaled.service
[Unit]
Description=Tailscale node agent
Documentation=https://tailscale.com/kb/
After=network.target
Wants=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/tailscaled
ExecStartPre=/bin/mkdir -p /var/lib/tailscale /run/tailscale
ExecStart=/bin/sh -c 'exec /usr/sbin/tailscaled --state=/var/lib/tailscale/tailscaled.state --socket=/run/tailscale/tailscaled.sock --port=41641 ${FLAGS}'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
